# ai-camera.service
[Unit]
Description=AI Camera Service
After=network.target

[Service]
Type=simple
User=dennisl
WorkingDirectory=/home/dennisl/AI-Camera
ExecStart=/home/dennisl/AI-Camera/start.sh
Restart=always
RestartSec=10
StandardOutput=append:/home/dennisl/AI-Camera/logs/stdout.log
StandardError=append:/home/dennisl/AI-Camera/logs/stderr.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
